# system/ssh_tls_proxy/sshd443.socket
[Unit]
Description=SSH-over-TLS Socket
PartOf=sshd443.service

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target
